# ansible/roles/media/files/mnt-media.automount
[Unit]
Description=Automount for the media library

[Automount]
Where=/mnt/media

[Install]
WantedBy=multi-user.target
